# cleanup-images.timer
[Unit]
Description=cleanup-images
Requires=docker.service

[Timer]
OnCalendar=daily

[X-Fleet]
Global=true
